# backup/com.user.restic-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.restic-backup</string>
    
    <key>ProgramArguments</key>
    <array>
        <!-- IMPORTANT: Update this path to match your setup -->
        <string>/Users/YOUR_USERNAME/machine-setup/backup/backup.sh</string>
    </array>
    
    <!-- Run daily at 2 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/restic-backup.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/restic-backup.log</string>
    
    <!-- Run if missed (e.g., computer was asleep) -->
    <key>StartInterval</key>
    <integer>86400</integer>
</dict>
</plist>
